# macworker/com.portfolioagent.macworker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.portfolioagent.macworker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/CHANGE_ME/PersonalPortfolioAnalyst/macworker/worker.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PA_API_BASE</key>
        <string>http://localhost:8000</string>
        <key>PA_API_TOKEN</key>
        <string>CHANGE_ME_API_TOKEN</string>
        <key>PA_IMESSAGE_RECIPIENT</key>
        <string>+15555550123</string>
    </dict>

    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/portfolioagent.macworker.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/portfolioagent.macworker.err.log</string>
</dict>
</plist>
